# systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata
   
[Service]
Type=simple
ExecStart=/usr/sbin/kanata --cfg %h/.config/kanata/kanata.kbd
Restart=on-failure
   
[Install]
WantedBy=default.target
